# etc/af-monalisa.cron
* * * * * sysman "/opt/aaf/bin/af-monalisa.pl" > /dev/null 2> /dev/null
